# better-chatbot.service
[Unit]
Description=Better Chatbot Application
After=network.target

[Service]
Type=simple
User=ric
WorkingDirectory=/home/ric/Projects/better-chatbot
Environment=PATH=/home/ric/.nvm/versions/node/v22.17.0/bin:/home/ric/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NO_HTTPS=1
ExecStart=/home/ric/.local/share/pnpm/pm2-runtime start ecosystem.config.cjs
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target